# templates/homeassistant.service.j2
[Unit]
Description=Home Assistant
After=syslog.target network.target

[Service]
ExecStart=/opt/homeassistant/bin/hass
User=homeassistant

[Install]
WantedBy=multi-user.target